# homecore-tools/rootfs/etc/services.d/hct-daemon/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finalização do serviço HCT Daemon
# ==============================================================================

bashio::log.info "Serviço HCT Daemon finalizado"
