# homecore-tools/rootfs/etc/services.d/hct-daemon/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Serviço HCT Daemon (S6 Overlay)
# ==============================================================================

bashio::log.info "Iniciando serviço HCT Daemon..."
exec /usr/bin/hct_daemon.py
