# homecore_tools/rootfs/etc/s6-overlay/s6-rc.d/homecore-tools/run
#!/usr/bin/execlineb -P
with-contenv
exec /app/run.sh
